# etc/blecastd.service
[Unit]
Description=blecastd BLE Advertising Beacon Daemon
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=root
Group=root
Environment=PYTHONPATH=/opt/blecastd/src
ExecStart=/usr/bin/python3 -m blecastd.main --config /etc/blecastd/blecastd.toml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
